# terrarium.service
[Unit]
Description=Terrarium Monitor
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi/terrarium
ExecStart=/home/pi/terrarium/venv/bin/python main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
